# config/systemd/user/waybar.service.d/restart.conf
# Make the restart loop robust: no rate-limit, 1s backoff.
# This replaces the old start-waybar.sh retry logic.
[Unit]
StartLimitIntervalSec=0

[Service]
RestartSec=1
